# CustomFiles/Depends/rc.local
export PATH='/opt/bin:/opt/sbin:/opt/bin:/opt/sbin:/usr/sbin:/usr/bin:/sbin:/bin'
cd /opt/.aria2/ && sh tracker.sh
. /opt/etc/init.d/rc.unslung start
exit 0 